# scripts/lootbox-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lootbox.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/johnlindquist/.bun/bin/bun</string>
        <string>run</string>
        <string>/Users/johnlindquist/dev/lootbox-bun/src/lib/observability/server-monitor.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/johnlindquist/dev/lootbox-bun</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/johnlindquist/.bun/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/johnlindquist/.lootbox-server-logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/johnlindquist/.lootbox-server-logs/launchd-stderr.log</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>8192</integer>
    </dict>
</dict>
</plist>